# build/supervisord.conf
[supervisord]
nodaemon=true

[program:backend]
priority=1
command=telemetry
startsecs=5
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:frontend]
priority=2
command=start-frontend.sh
startsecs=5
stopwaitsecs=10
autorestart=true
directory=/usr/src/app/
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0